# root/etc/systemd/system/bootc-fetch-apply-updates.service.d/dedupe-bound-images.conf
# bootc's bound-image prune fails hard if a stale image ID carries more than
# one tag (identical builds published under two versions). Untag duplicates
# first so the nightly upgrade can't wedge on it.
[Service]
ExecStartPre=/usr/local/sbin/dedupe-bound-images.sh
