# root/usr/share/containers/systemd/rclone-move-media/rclone-move-media.container
[Unit]
Description=rclone-move-media (uploads transcoded media, refreshes plex)
After=download-secrets.service
Requires=download-secrets.service
RequiresMountsFor=/var/mnt/data

[Container]
Image=ghcr.io/jonohill/rclone-move:148
GlobalArgs=--storage-opt=additionalimagestore=/usr/lib/bootc/storage
# Image is logically bound, so bootc keeps it in the additional store above;
# never pull into the primary store (avoids duplicate copies / orphaned records).
Pull=never
ContainerName=rclone_move_media
Network=media.network
Volume=/mnt/data/tmp/plex-media:/source
Volume=/mnt/data/rclone_move_media:/config
Environment=SOURCE=/source
Environment=DEST=jotta:/plex-media
Environment=PLEX_PREFIX=/data/plex-media
Environment=PLEX_URL=http://plex:32400
# PLEX_TOKEN=...
EnvironmentFile=/run/rclone-move-media/env
Environment=RCLONE_EXTRA_FLAGS=--exclude,*.part.*
Exec=--metrics-addr :9090

[Service]
Restart=always

[Install]
WantedBy=multi-user.target
