# root/usr/share/containers/systemd/hermes-agent/hermes-dashboard.container
[Unit]
Description=hermes-agent dashboard
RequiresMountsFor=/var/mnt/data
After=hermes-agent.service

[Container]
Image=docker.io/nousresearch/hermes-agent:v2026.6.19
GlobalArgs=--storage-opt=additionalimagestore=/usr/lib/bootc/storage
Pod=hermes-agent.pod
Environment=TZ=Pacific/Auckland
Environment=HERMES_UID=1000
Environment=HERMES_GID=1000
Volume=/mnt/data/hermes:/opt/data
# Required to disable host check, have websockets work
# Fine because it's on the pod's network
Exec=dashboard --host 0.0.0.0 --insecure --no-open

[Service]
Restart=always
# Ensure the bind-mount source exists on the data volume before podman runs;
# RequiresMountsFor (above) guarantees /var/mnt/data is mounted first.
ExecStartPre=/usr/bin/mkdir -p /mnt/data/hermes

[Install]
WantedBy=multi-user.target
