# root/usr/share/containers/systemd/cutout/whisper-batch.container
[Unit]
Description=cutout whisper-batch (transcription sidecar)
RequiresMountsFor=/var/mnt/data

[Container]
Image=ghcr.io/jonohill/whisper-batch:20
GlobalArgs=--storage-opt=additionalimagestore=/usr/lib/bootc/storage
# Image is logically bound, so bootc keeps it in the additional store above;
# never pull into the primary store (avoids duplicate copies / orphaned records).
Pull=never
Pod=cutout.pod
ContainerName=whisper-batch
Environment=TZ=Pacific/Auckland
Environment=WHISPER_MODEL_NAME=base.en-q8_0
# Downloaded model weights persist here so they survive container restarts.
# :U makes podman chown the source to the container user (app, uid 10001) so
# the model can be written -- the ExecStartPre mkdir creates it as root.
Volume=/mnt/data/cutout/whisper-models:/models:U

[Service]
Restart=always
# Lower CPU scheduling priority so transcription yields to other work.
Nice=5
# Ensure the bind-mount source exists on the data volume before podman runs;
# RequiresMountsFor (above) guarantees /var/mnt/data is mounted first.
ExecStartPre=/usr/bin/mkdir -p /mnt/data/cutout/whisper-models

[Install]
WantedBy=multi-user.target
