# root/usr/share/containers/systemd/qbittorrent/qbittorrent-cloudflared.container
[Unit]
Description=qbittorrent cloudflared tunnel
After=download-secrets.service qbittorrent.service
Requires=download-secrets.service
Wants=qbittorrent.service

[Container]
Image=docker.io/cloudflare/cloudflared:2026.7.0
GlobalArgs=--storage-opt=additionalimagestore=/usr/lib/bootc/storage
# Image is logically bound, so bootc keeps it in the additional store above;
# never pull into the primary store (avoids duplicate copies / orphaned records).
Pull=never
Network=qbt.network
Environment=TZ=Pacific/Auckland
Volume=/run/qbittorrent-cloudflared/config.yml:/etc/cloudflared/config.yml:ro
Volume=/run/qbittorrent-cloudflared/credentials.json:/etc/cloudflared/credentials.json:ro
SecurityLabelDisable=true
Exec=tunnel run

[Service]
Restart=always
EnvironmentFile=/run/qbittorrent-cloudflared/env
ExecStartPre=/usr/bin/sh -c 'cat /run/cloudflared-access/access.yml /run/qbittorrent-cloudflared/tunnel.yml /usr/share/containers/systemd/qbittorrent/config.yml > /run/qbittorrent-cloudflared/config.yml'
ExecStartPre=/usr/bin/sh -c 'echo "$$TUNNEL_CRED_CONTENTS" > /run/qbittorrent-cloudflared/credentials.json'

[Install]
WantedBy=multi-user.target
