# root/usr/share/containers/systemd/goose/goose.container
[Unit]
Description=goose Slack agent (OpenAB)
After=download-secrets.service
Requires=download-secrets.service
RequiresMountsFor=/var/mnt/data

[Container]
Image=ghcr.io/jonohill/goose-claw:5
GlobalArgs=--storage-opt=additionalimagestore=/usr/lib/bootc/storage
# Image is logically bound, so bootc keeps it in the additional store above;
# never pull into the primary store (avoids duplicate copies / orphaned records).
Pull=never
Pod=goose.pod
ContainerName=goose
Environment=TZ=Pacific/Auckland
# Read-only host journal (= container logs; the journald log driver is used)
# via systemd-journal-gatewayd on the goose bridge gateway. e.g.
#   curl -H 'Accept: application/json' -H 'Range: entries=:-100:100' \
#     "$JOURNAL_GATEWAY_URL/entries?CONTAINER_NAME=plex"
Environment=JOURNAL_GATEWAY_URL=http://10.89.3.1:19531
# Slack tokens (xoxb bot + xapp app-level) and the provider key. OpenAB expands
# ${SLACK_BOT_TOKEN}/${SLACK_APP_TOKEN}/${OLLAMA_API_KEY} from its environment
# when it loads the baked /etc/openab/config.toml, then forwards the resolved
# provider values to each `goose acp` child. Required keys:
#   SLACK_BOT_TOKEN, SLACK_APP_TOKEN, OLLAMA_API_KEY
EnvironmentFile=/run/goose/env
# Persist config + sessions + OpenAB state (~/.config/goose, ~/.openab,
# /data/workspace). :U chowns the source to the container user (goose, uid 1000).
Volume=/mnt/data/goose:/data:U
Volume=/mnt/data/downloads:/mnt/data/downloads
# No Exec= override: use the image default CMD
# (openab run -c /etc/openab/config.toml). Slack Socket Mode is an outbound
# WebSocket, so no ports are published and no inbound proxy/tunnel is needed.

[Service]
Restart=always
# Ensure the bind-mount source exists on the data volume before podman runs;
# RequiresMountsFor (above) guarantees /var/mnt/data is mounted first.
ExecStartPre=/usr/bin/mkdir -p /mnt/data/goose

[Install]
WantedBy=multi-user.target
