# root/usr/share/containers/systemd/goose/goose-cloudflared.container
[Unit]
Description=goose cloudflared tunnel
After=download-secrets.service
Requires=download-secrets.service

[Container]
Image=docker.io/cloudflare/cloudflared:2026.7.0
GlobalArgs=--storage-opt=additionalimagestore=/usr/lib/bootc/storage
# Image is logically bound, so bootc keeps it in the additional store above;
# never pull into the primary store (avoids duplicate copies / orphaned records).
Pull=never
Pod=goose.pod
Environment=TZ=Pacific/Auckland
Volume=/run/goose-cloudflared/config.yml:/etc/cloudflared/config.yml:ro
Volume=/run/goose-cloudflared/credentials.json:/etc/cloudflared/credentials.json:ro
SecurityLabelDisable=true
Exec=tunnel run

[Service]
Restart=always
EnvironmentFile=/run/goose-cloudflared/env
ExecStartPre=/usr/bin/sh -c 'cat /run/cloudflared-access/access.yml /run/goose-cloudflared/tunnel.yml /usr/share/containers/systemd/goose/config.yml > /run/goose-cloudflared/config.yml'
ExecStartPre=/usr/bin/sh -c 'echo "$$TUNNEL_CRED_CONTENTS" > /run/goose-cloudflared/credentials.json'

# Disabled for now — not auto-started at boot.
# [Install]
# WantedBy=multi-user.target
