# root/usr/share/containers/systemd/goose/goose-api-proxy.container
[Unit]
Description=goose API reverse proxy
After=goose.service
# Consumes the external gate secret decrypted to /run/goose/caddy.env.
After=download-secrets.service
Requires=download-secrets.service

[Container]
Image=docker.io/library/caddy:2-alpine
GlobalArgs=--storage-opt=additionalimagestore=/usr/lib/bootc/storage
# Image is logically bound, so bootc keeps it in the additional store above;
# never pull into the primary store (avoids duplicate copies / orphaned records).
Pull=never
Pod=goose.pod
Environment=TZ=Pacific/Auckland
# Internal secret (minted by goose.pod): injected upstream as X-Secret-Key so
# goosed accepts the request. Referenced as {env.GOOSE_SERVER__SECRET_KEY}.
EnvironmentFile=/run/goose/secret.env
# External secret: the value clients must present in X-Secret-Key. Referenced in
# the Caddyfile as {$GOOSE_EXTERNAL_SECRET_KEY} (substituted at config load).
EnvironmentFile=/run/goose/caddy.env
Volume=/usr/share/containers/systemd/goose/Caddyfile:/etc/caddy/Caddyfile:ro

[Service]
Restart=always
# Fail closed: ensure secret key is set
ExecStartPre=/usr/bin/sh -c 'grep -Eq "^GOOSE_EXTERNAL_SECRET_KEY=.+" /run/goose/caddy.env'

[Install]
WantedBy=multi-user.target
