# root/usr/share/containers/systemd/plex/plex.container
[Unit]
Description=plex
After=download-secrets.service
Requires=download-secrets.service
RequiresMountsFor=/var/mnt/data

[Container]
Image=ghcr.io/hillnz/docker-nas-plex:528
GlobalArgs=--storage-opt=additionalimagestore=/usr/lib/bootc/storage
# Image is logically bound, so bootc keeps it in the additional store above;
# never pull into the primary store (avoids duplicate copies / orphaned records).
Pull=never
Pod=plex.pod
ContainerName=plex
Environment=TZ=Pacific/Auckland
Environment=HOSTNAME=larry
Environment=INCOMING_DIR=/data/incoming
Environment=RCLONE_MOUNT_DIR=jotta:/plex-media
Environment=RCLONE_MOUNT_CACHE_DIR=/caches/rclone_plex
Environment=RCLONE_MOUNT_TARGET=/data/plex-media
Environment=RCLONE_CACHE_MAX_SIZE=25G
Environment=PROXY_PORT=36000
Environment=PROXY_PLEX_LIBRARY_PATH=/data/plex-media
Environment=RUST_LOG=info
Environment=PROXY_TARGET=https://proxy-worker.jhill.workers.dev
Environment=TS_HOSTNAME=plex
Environment=TS_DRIVE_SHARES=media:/data/plex-media
Environment=FILESTASH_BASE=/files
EnvironmentFile=/run/plex/env
# AddHost lives on the pod (plex.pod) -- podman 5 rejects network options on a
# container that shares a pod's network namespace.
Volume=/mnt/data/plex:/config
Volume=/mnt/data/tmp/plex-transcode:/transcode
Volume=/tmp:/tmp
Volume=/var/mnt/data/cache/rclone-plex:/caches/rclone_plex
Volume=/mnt/data/downloads:/downloads
AddDevice=/dev/fuse
AddCapability=SYS_ADMIN
# For FUSE
SecurityLabelDisable=true
HealthStartPeriod=120s

[Service]
Restart=always

[Install]
WantedBy=multi-user.target
