# root/usr/share/containers/systemd/av1ator/av1ator.container
[Unit]
Description=av1ator (AV1 transcoder)
RequiresMountsFor=/var/mnt/data

[Container]
Image=ghcr.io/jonohill/av1ator:5
GlobalArgs=--storage-opt=additionalimagestore=/usr/lib/bootc/storage
# Image is logically bound, so bootc keeps it in the additional store above;
# never pull into the primary store (avoids duplicate copies / orphaned records).
Pull=never
ContainerName=av1ator
Volume=/mnt/data/tmp/sorted-media:/input
Volume=/mnt/data/tmp/plex-media:/output
Exec=--crf 27 --preset 6 --delete-input --min-size 500MB --max-width 1920 --skip-codec av1 --passthrough-filtered --nice 10

[Service]
Restart=always

[Install]
WantedBy=multi-user.target
