# root/etc/systemd/system/systemd-journal-gatewayd.socket.d/goose-bridge.conf
[Socket]
# Serve the journal (read-only, by design of gatewayd) to the goose agent.
# Bind only to the goose bridge gateway (pinned in goose.network) so it is
# never exposed on the VNIC/tailscale interfaces. Note it has no auth, so any
# container on any local bridge can technically route to it via the host.
ListenStream=
ListenStream=10.89.3.1:19531
# The bridge interface only exists once the goose network is up; FreeBind lets
# the socket bind at boot regardless of ordering.
FreeBind=true
