# root/etc/systemd/system/bootc-fetch-apply-updates.timer.d/schedule.conf
# Pin bootc's auto-update check to a predictable window instead of the base
# image's randomised daily schedule. The .service runs `bootc upgrade --apply`,
# which only reboots when a newer image was actually pulled.
[Timer]
# Clear every schedule key the base unit might set before defining our own.
OnBootSec=
OnStartupSec=
OnActiveSec=
OnUnitActiveSec=
OnUnitInactiveSec=
OnCalendar=
OnCalendar=*-*-* 03:00:00 Pacific/Auckland
# Fire at exactly 03:00 (the base unit spreads the fleet with a random delay;
# this is a single host, so keep it deterministic).
RandomizedDelaySec=0
# If the host was powered off at 03:00, run the check at next boot.
Persistent=true
