# root/etc/systemd/system/bootc-fetch-apply-updates.timer.d/schedule.conf
# Pin bootc's auto-update check to a predictable window instead of the base
# image's randomised daily schedule. The .service runs `bootc upgrade --apply`,
# which only reboots when a newer image was actually pulled.
[Timer]
# Clear every schedule key the base unit might set before defining our own.
OnBootSec=
OnStartupSec=
OnActiveSec=
OnUnitActiveSec=
OnUnitInactiveSec=
OnCalendar=
OnCalendar=*-*-* 05:00:00 Pacific/Auckland
RandomizedDelaySec=0
Persistent=true
