# vendor/bootc-secrets/install/download-secrets.service
[Unit]
Description=Download and decrypt host secrets to /run
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/usr/local/sbin/download_secrets.sh

[Install]
WantedBy=multi-user.target

# Workloads that need the decrypted secrets should order themselves after this
# unit. For a Quadlet container, add to its [Unit] section:
#
#   [Unit]
#   After=download-secrets.service
#   Requires=download-secrets.service
#
# and consume the secrets via e.g. EnvironmentFile=/run/<service>/env.
